# heketi/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:heketi]
command=/usr/local/bin/start-heketi.sh 